# SCRIPTS/launchd/com.r2c.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.r2c.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/hkobayashi/projects/commerce-faq-tasks/SCRIPTS/r2c-cron-wrapper.sh</string>
        <string>--script</string>
        <string>monitor-claude-health.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/hkobayashi/projects/commerce-faq-tasks</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/hkobayashi</string>
        <key>CLAUDE_CONFIG_DIR</key>
        <string>/Users/hkobayashi/.claude-r2c-config</string>
    </dict>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/hkobayashi/.claude-r2c-config/logs/launchd-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hkobayashi/.claude-r2c-config/logs/launchd-monitor-err.log</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
