# SCRIPTS/launchd/com.r2c.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        R2C 24h Loop - Supervisor (1 分間隔)
        雛形配置: SCRIPTS/launchd/com.r2c.supervisor.plist (tracked)
        実体配置: ~/Library/LaunchAgents/com.r2c.supervisor.plist (hkobayashi が手動 cp + load)

          cp SCRIPTS/launchd/com.r2c.supervisor.plist ~/Library/LaunchAgents/
          launchctl load -w ~/Library/LaunchAgents/com.r2c.supervisor.plist

        役割: r2c-cron-wrapper.sh 経由で r2c-supervisor.sh を起動。
              running 状態の Lane を監視し、MAX_RUN_MINUTES=45 超過の stuck Lane を
              auto-retry / rollback する (1 分間隔ポーリング)。
    -->
    <key>Label</key>
    <string>com.r2c.supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/hkobayashi/Documents/GitHub/commerce-faq-tasks/SCRIPTS/r2c-cron-wrapper.sh</string>
        <string>--script</string>
        <string>r2c-supervisor.sh</string>
    </array>

    <!-- 1 分間隔 -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/hkobayashi/.claude-r2c-config/logs/launchd-supervisor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hkobayashi/.claude-r2c-config/logs/launchd-supervisor-err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/hkobayashi</string>
        <key>CLAUDE_CONFIG_DIR</key>
        <string>/Users/hkobayashi/.claude-r2c-config</string>
        <key>CLAUDE_CODE_EXPERIMENTAL_AGENT_TEAMS</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/hkobayashi/Documents/GitHub/commerce-faq-tasks</string>

    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
